# deployment_scripts/puppet/modules/contrail/templates/72-contrail-dpdk-on-vf.rules.erb
SUBSYSTEM=="net", ACTION=="add", KERNELS=="<%= @dpdk_dev_pci %>", NAME="<%= @dpdk_dev_name %>", RUN+="/bin/ip link set dev %k address <%= @dpdk_dev_mac %>"
