# deployment_scripts/puppet/modules/contrail/templates/contrail-tor-agent.ini.erb

[program:contrail-tor-agent-<%= @tor_id -%>]
command=/usr/bin/contrail-tor-agent --config_file /etc/contrail/contrail-tor-agent-<%= @tor_id -%>.conf
priority=420
autostart=true
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-tor-agent-<%= @tor_id -%>-stdout.log
stderr_logfile=/dev/null
startsecs=5
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)