# deployment_scripts/puppet/modules/contrail/templates/contrail-vrouter-dpdk.ini.erb
[program:contrail-vrouter-dpdk]
command=taskset <%= scope.lookupvar('contrail::vrouter_core_mask') %> /usr/bin/contrail-vrouter-dpdk --no-daemon
priority=410
loglevel=debug
autostart=true
killasgroup=true
stdout_capture_maxbytes=1MB
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-vrouter-dpdk-stdout.log
stderr_logfile=/var/log/contrail/contrail-vrouter-dpdk-stderr.log
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)